# etc/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn3 -w 4 launch_web:app
directory=/srv/pyapp
